# files/telnetd.socket
[Unit]
Description=Telnetd socket

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
