# files/telnetd@.service
[Unit]
Description=Telnet Server

[Service]
ExecStart=/usr/sbin/in.telnetd
StandardInput=socket
